# crontab
* * * * *  php -f /app/file0.php 
* * * * *  php -f /app/file1.php 


